# docker/crontab
* * * * * /usr/local/bin/php /var/www/html/cron/cron.php >> /var/log/cron.log 2>&1
